# moneybotsapp.service
[Unit]
Description=MoneyBots Go Application
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/moneybotsapp
ExecStart=/opt/moneybotsapp/moneybotsapp
Restart=always
RestartSec=5
Environment=PORT=8080
Environment=APP_VERSION=1.0.0

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/moneybotsapp

[Install]
WantedBy=multi-user.target